# roles/mastodon/templates/mastodon-clean.timer.j2
[Unit]
Description=Run mastodon clean script daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=43m
Persistent=true

[Install]
WantedBy=timers.target
