# roles/backup/templates/backup.service.j2
[Unit]
Description=Backup
OnFailure=status-email-root@backup.service

[Service]
Type=oneshot
ExecStart="{{ root_bin }}backup.sh
PrivateTmp=true
